# devel/etc/catkin/profile.d/skiros.bash
/home/krishna/Downloads/final_pkg_ws/devel/.private/skiros2/etc/catkin/profile.d/skiros.bash